# systemd/.config/systemd/user/daily-wallpaper.service
[Unit]
Description=Download and apply the daily wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/daily-wallpaper.sh
